# tools/docker-files/pencmd/target/usr/lib/systemd/system/pen-cmd.service
[Unit]
Description=Pensando Cluster Manager Daemon
Documentation=http://pensando.io/
StopWhenUnneeded=yes
After=network.target

[Service]
EnvironmentFile=-/etc/pensando/pencmd.conf
ExecStartPre=/usr/bin/ln -s /docker.sock /var/run
ExecStart=/usr/bin/pen-cmd

Type=simple
Restart=always
RestartSec=10s
KillMode=process

[Install]
WantedBy=pensando.target
