# tools/docker-files/pencmd/target/usr/lib/systemd/system/pen-apigw.service
[Unit]
Description=Pensando API Gateway

[Service]
ExecStartPre=-/usr/bin/docker kill pen-apigw
ExecStartPre=-/usr/bin/docker rm -f pen-apigw
ExecStart=/usr/pensando/bin/systemd-docker \
	 --cgroups name=systemd --cgroups=cpu  run --rm --net=host --name pen-apigw \
	-v /etc/pensando:/etc/pensando \
	pen-apigw
ExecStop=/usr/bin/docker stop -t 2 pen-apigw
ExecStopPost=/usr/bin/docker rm -f pen-apigw
Restart=always
RestartSec=10s
Type=notify
NotifyAccess=all
TimeoutStartSec=120
TimeoutStopSec=15
