# tools/docker-files/cmd/target/usr/lib/systemd/system/pen-kubelet.service
[Unit]
Description=Kubelet

[Service]
User=root
EnvironmentFile=-/etc/pensando/kubelet.conf
ExecStart=/usr/pensando/bin/kubelet \
            --node-status-update-frequency 4s \
			$KUBELET_ADDRESS \
			$ALLOW_PRIVILEGED \
			$KUBE_CONFIG \
			$REQUIRE_KUBE_CONFIG \
			$NODEIP \
			$HOSTNAME_OVERRIDE \
			$CGROUP_DRIVER 
Restart=always
RestartSec=10s
LimitNOFILE=40000
