# utils/systemd/test/pensando.target
[Unit]
Description=desc_target
DefaultDependencies=no
