# utils/systemd/test/service1.service
[Unit]
Description=desc_service1
StopWhenUnneeded=yes
DefaultDependencies=no

[Install]
WantedBy=pensando.target

[Service]
ExecStart=/bin/sleep 400
